# systemd/dialpad-webhook.service
[Unit]
Description=Dialpad SMS Webhook Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/projects/skills/work/dialpad
EnvironmentFile=%h/projects/skills/work/dialpad/.env
# NOTE: Adjust Python path to your system (e.g., /usr/bin/python3, ~/.nix-profile/bin/python3)
ExecStart=%h/.linuxbrew/bin/python3 %h/projects/skills/work/dialpad/webhook_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target